# scripts/crontab
5 1 * * * cd /opt/audiodescription && make cron-import >> /opt/audiodescription/logs/cron-import.log